# crontab
*/1 * * * * cd /data && node ticker.js > ${STDOUT_LOC} 2> ${STDERR_LOC}


